# docker/supervisor/conf.d/ldapproxy.sv.conf
[program:ldapproxy]
command=/bin/bash /opt/ldapproxy.sh
redirect_stderr=true
autostart=true
autorestart=false