# docker/supervisor/conf.d/netboot.sv.conf
[program:netboot]
command=/bin/bash /opt/netboot.sh
redirect_stderr=true
autostart=true
autorestart=false