# docker/supervisor/conf.d/webapp.sv.conf
[program:webapp]
command=/bin/bash /opt/webapp.sh
redirect_stderr=true
autostart=true
autorestart=true
